# docker/8.0/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php]
<<<<<<< HEAD
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan serve --host=0.0.0.0 --port=80
=======
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/PrivyUpdate/artisan serve --host=143.110.241.48 --port=80
>>>>>>> 6a58521c61b190331bb05066a19296515dee7a5a
user=sail
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
